# services/system/beo-librespot.service
[Unit]
Description=BeoSound5c go-librespot (Spotify Connect local playback)
After=network-online.target sound.target
Wants=network-online.target
Before=beo-player-local.service
OnFailure=beo-notify-failure@%n.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=/etc/beosound5c
Environment=XDG_RUNTIME_DIR=/run/user/__UID__
ExecStart=/usr/local/bin/go-librespot --config_dir /etc/beosound5c/librespot
# Wait for the API to come up, but BOUNDED (30s): an unbounded loop hangs
# until TimeoutStartSec (90s default) on bad credentials. Failing here marks
# the unit failed deliberately — Restart=on-failure + StartLimit take over.
ExecStartPost=/bin/sh -c 'timeout 30 sh -c "until curl -sf http://localhost:3678/status >/dev/null 2>&1; do sleep 0.2; done" || { echo "go-librespot API not reachable after 30s" >&2; exit 1; }'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
